# templates/systemd/autostart.service.j2
[Unit]
Description=Autostart configuration for {{ item }}
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart={{ item }}

[Install]
WantedBy=sway-session.target
